# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:squash-tm]
command=/usr/share/squash-tm/bin/startup.sh